# systemd/shield.timer
[Unit]
Description=PAM Shield Purge Routine

[Timer]
OnCalendar=daily
Persistent=false

[Install]
WantedBy=timers.target
